# containers/weft_time/config/supervisord.conf
[supervisord]
nodaemon=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
pidfile=/run/supervisord.pid
user=root

[program:chrony]
command=chronyd -f /etc/chrony/chrony.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=true
startretries=20
priority=200

[program:gpsd]
command=gpsd -N -G -n /dev/ttyACM0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=true
priority=1


